# tools/systemd/open-intranet.target
[Unit]
Description="OpenIntranet instance"
#Requires=

[Install]
WantedBy=multi-user.target
